# packages/evolve-agent/evolve_agent-1.0.2-py2-none-any.whl/etc/evolve-agent/upstart/evolve-agent.conf
description "Evolve Agent"

start on runlevel [2345]
stop on runlevel [!2345]

script
    exec /usr/bin/python /usr/bin/evolve-agent -p /var/run/evolve-agent.pid
end script
