# packages/fb2feed/fb2feed-0.3.tar.gz/fb2feed-0.3/systemd/fb2feed-timer.target
[Unit]
Description=Fb2Feed Timer Target
StopWhenUnneeded=yes