# packages/evolve-agent/evolve_agent-1.0.1-py2-none-any.whl/etc/evolve-agent/systemd/evolve-agent.service
[Unit]
Description=evolve-agent

[Service]
ExecStart=/usr/bin/python /usr/bin/evolve-agent -d
Type=forking
PIDFile=/var/run/evolve-agent.pid

[Install]
WantedBy=multi-user.target
